# files/nvidia/nvidia-eula.service
[Unit]
Description=Display Nvidia driver EULA

[Service]
Type=oneshot
ExecStart=/etc/eks/nvidia-eula

[Install]
WantedBy=basic.target